# templates/systemd-nspawn.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-nspawn --quiet --keep-unit --boot --machine=%I
